# deploy/twilight.service
[Unit]
Description=Twilight Go API
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/root/Twilight
ExecStart=/root/Twilight/bin/twilight api --host 127.0.0.1 --port 5000 --config /root/Twilight/config.toml

Environment=TWILIGHT_CONFIG_FILE=/root/Twilight/config.toml
EnvironmentFile=-/root/Twilight/.env

LimitNOFILE=65535
LimitNPROC=4096
MemoryMax=1G
MemoryHigh=768M

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=twilight

[Install]
WantedBy=multi-user.target
